# lambton/2020/fall/ese3005/class_2/myservice.service
[Unit]
Description=a test service experiment!

[Service]
Type=simple
ExecStart=/usr/bin/test_service

[Install]
WantedBy=multi-user.target
